# multinode-rails/rails/etc/service/web_app/run
#!/bin/bash

cd /web_app

passenger start
