# gitlab-ci-runner/etc/service/gitlab-ci-runner/run
#!/bin/bash

RUN cd gitlab-ci-runner.git
exec bundle exec ./bin/setup_and_run
